# linux/crontab/chezmoi
0 * * * * $HOME/.local/bin/chezmoi -v apply
