# linux/crontab/chezmoi
0 * * * * ~/.local/bin/chezmoi -v apply
